# devops/backup.cron
1 22 * * * __USER__ __CABINET_PATH__/devops/backup.sh >> __CABINET_PATH__/devops/cron.log 2>&1

