# devops/backup.cron
1 22 * * * oasis /home/oasis/hwcentral-cabinet/devops/backup.sh >> /home/oasis/hwcentral-cabinet/devops/cron.log 2>&1
